# scrappers/crontab
#run scrappers every minute
* 3 * * * /usr/local/bin/python -u /app/src/main.py > /proc/1/fd/1 2>/proc/1/fd/2
